# ray/systemd/ray-head.service
[Unit]
Description=Ray Head Node (vLLM cluster)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=soypete
WorkingDirectory=/home/soypete
ExecStartPre=/bin/bash -c 'ip addr show enp1s0f0np0 | grep -q 192.168.100.10 || sudo ip addr add 192.168.100.10/24 dev enp1s0f0np0'
ExecStart=/bin/bash /home/soypete/start-head.sh
Restart=on-failure
RestartSec=30
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
